# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log

[program:sparsh]
command=node server.js
directory=/app/sparsh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sparsh.log
stderr_logfile=/var/log/supervisor/sparsh.log
environment=PORT=3000

[program:hamraaz]
command=node server.js
directory=/app/hamraaz
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/hamraaz.log
stderr_logfile=/var/log/supervisor/hamraaz.log
environment=PORT=3001

[program:sbi]
command=node server.js
directory=/app/sbi
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sbi.log
stderr_logfile=/var/log/supervisor/sbi.log
environment=PORT=3002